# deploy/pf-scout-api.service
[Unit]
Description=pf-scout API server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pf-scout-bot/scout-api
EnvironmentFile=/opt/pf-scout-bot/.env
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8420
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
